# contrib/arsenal.service
[Unit]
Description=Arsenal team server
After=network.target

[Service]
ExecStart=/bin/bash -c "cd /opt/arsenal/teamserver; source /opt/arsenal/venv/bin/activate; uwsgi --ini uwsgi.ini --lazy-apps"
ExecReload=/bin/kill -SIGQUIT $MAINPID
KillSignal=SIGTERM
StandardError=syslog
Type=simple

[Install]
WantedBy=multi-user.target